# modules/statsd/files/upstart.statsd.conf
description "statsd"

start on startup
stop on shutdown

respawn

script
    exec sudo -u statsd /opt/statsd/start.sh
end script

